# supervisord.conf
## service_script.conf
[supervisord]
user=root
nodaemon=true

[program:fuse]
user=root
command=python myfuse.py /mailbox

[program:nginx]
command=./entrypoint.sh && nginx -g "daemon off;"